# http_dyndns.timer
[Unit]
Description=Run the Http DynDNS service periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
