# http_dyndns.service
[Unit]
Description=HTTP DynDNS Updater

[Service]
Type=simple
ExecStart=/usr/local/bin/HttpDynDNS.sh

[Install]
WantedBy=default.target